# systemd/system/mitmweb.service
[Unit]
Description=mitmproxy Web UI (tun mode)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mitmproxy
Group=mitmproxy
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW

; settings
Environment=MITM_TUN=mitm-tun
Environment=MITM_WEB_HOST=10.64.8.64
Environment=MITM_WEB_PORT=8081
Environment=MITM_WEB_PASSWORD=securitytheater

; unprivileged user creates tun before mitmweb starts
ExecStartPre=/bin/bash -c '/usr/sbin/ip link del ${MITM_TUN} 2>/dev/null || true'
ExecStartPre=/usr/sbin/ip tuntap add dev ${MITM_TUN} mode tun
ExecStartPre=/usr/sbin/ip link set ${MITM_TUN} up

; start mitmweb in tun mode
ExecStart=/usr/bin/mitmweb -vv --mode tun:${MITM_TUN} --web-host ${MITM_WEB_HOST} --web-port ${MITM_WEB_PORT} --set web_password=${MITM_WEB_PASSWORD}

; wait tun interface appears, then set up routing table to capture dev nat64mitm
ExecStartPost=/bin/bash -c 'for i in $(seq 1 60); do ip link show ${MITM_TUN} && break; sleep 1; done'
ExecStartPost=/usr/sbin/ip link set ${MITM_TUN} up
ExecStartPost=/usr/sbin/ip rule add fwmark 8080 lookup 8080
ExecStartPost=/usr/sbin/ip route replace default dev ${MITM_TUN} table 8080

Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
